# platforms/x86_64-accton_as9516_32d-r0/overlay/etc/systemd/system/bf_fpga_module.service
[Unit]
Description=bf_fpga Module
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe bf_fpga intr_mode="none"
ExecStop=/sbin/rmmod bf_fpga

[Install]
WantedBy=multi-user.target
