# platforms/x86_64-accton_as9516_32d-r0/overlay/etc/systemd/system/bf_knet_module.service
[Unit]
Description=bf_knet Module
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe bf_knet intr_mode="none"
ExecStop=/sbin/rmmod bf_knet

[Install]
WantedBy=multi-user.target
