# platforms/x86_64-accton_as9516_32d-r0/overlay/etc/systemd/system/bf_kdrv_module.service
[Unit]
Description=bf_kdrv Module
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe bf_kdrv intr_mode="msi"
ExecStop=/sbin/rmmod bf_kdrv

[Install]
WantedBy=multi-user.target
